# services/user/hewo_face.service.template
[Unit]
Description=HeWo Face node
Requires=xhost.service


[Service]
Type=simple
ExecStart=/home/__USER__/HeWo/run_module.sh face

[Install]
WantedBy=default.target
